# project-bi-engine-service/docker/container_files/etc/supervisord.d/project-bi-engine.ini
[program:project-bi-engine]
command = java APP_JAVA_OPTS -javaagent:/opt/apm/pinpoint-bootstrap-1.7.1.jar -Dpinpoint.agentId=project-bi-engine-Service -Dpinpoint.applicationName=project-bi-engine-App -jar /opt/apps/project-bi-engine.jar
stderr_logfile = syslog
stdout_logfile = syslog